# test_monitor.service
[Unit]
Description=Мониторинг процесса test

[Service]
Type=oneshot
ExecStart=/usr/local/bin/test_monitor.sh

